# www/bitrix/crontab/crontab.cfg

50 20 * * * run-one /usr/bin/php71 -f /home/b2153/domains/kirgu.ru/www/bitrix/php_interface/include/catalog_import/cron_frame.php 1 >/home/b2153/domains/kirgu.ru/www/bitrix/php_interface/include/catalog_import/logs/1.txt
